# localtime.service.in
[Unit]
Description=Timezone Updater Daemon

[Service]
User=USER
Group=GROUP

Nice=10
PrivateTmp=true
PrivateDevices=true
PrivateNetwork=true
ProtectSystem=full
ProtectControlGroups=true
ProtectKernelTunables=true
ProtectHome=true
NoNewPrivileges=true

ExecStart=BINDIR/localtimed

[Install]
WantedBy=default.target
