# localtime.service.in
[Unit]
Description=Timezone Updater Daemon
PartOf=geoclue.service

[Service]
User=USER
DynamicUser=true

Nice=10
PrivateTmp=true
PrivateDevices=true
PrivateNetwork=true
ProtectSystem=strict
ProtectControlGroups=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectHome=true
NoNewPrivileges=true
ProtectClock=true
ProtectKernelLogs=true
ProtectHostname=true
ProtectProc=noaccess
MemoryDenyWriteExecute=true

ExecStart=BINDIR/localtimed

[Install]
WantedBy=default.target
